# etc/profile.d/Makefile.am
etcdir = @sysconfdir@/profile.d
etc_DATA = Z97-byobu.sh
