# etc/profile.d/Makefile.am
etcdir = $(datadir)/@PACKAGE@/profile.d
etc_DATA = Z97-byobu.sh
