# etc/update-motd.d/Makefile.am
etcdir = $(datadir)/update-motd.d
etc_DATA = 55-window-manager
